# .build/supervisord.conf
[unix_http_server]
file=/root/supervisor.sock           ; path to your socket file
username= root
password=

[supervisord]
logfile=/root/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///root/supervisorctl.sock         ; use a unix:// URL  for a unix socket

[program:zk]
command=/bin/bash -c "/opt/kafka/bin/zookeeper-server-start.sh /opt/kafka/config/zookeeper.properties"
priority=5
autostart=true
startsecs=1
startretries=1000
autorestart=true
stopwaitseconds=1
stdout_logfile=/root/zookeeper.log
stderr_logfile=/root/zookeeper_error.log


[program:kafka]
command=/bin/bash -c "/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties"
priority=6
autostart=true
startsecs=1
startretries=1000
autorestart=true
stopwaitseconds=5
stdout_logfile=/dev/stdout   ; when directing to /dev/stdout, must have maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0